# distro/rpath-tools.cron
#m      h   DoM mon DoW     user    cmd
*       2   *   *   *       root    /usr/bin/rpath-tools register --random-wait
