# distro/rpath-tools.cron
#m      h   DoM mon DoW     user    cmd
10      2   *   *   *       root    /usr/bin/rpath-register --quiet --random-wait
22      3   *   *   *       root    /usr/bin/rpath-tmpwatch --quiet delete
